# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

; ---------------------------------------------------------------------------
; 1. Litestream – starts first, restores DB if needed, then keeps replicating
; ---------------------------------------------------------------------------
[program:litestream]
command=/app/scripts/litestream_start.sh
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=3
startsecs=5

; ---------------------------------------------------------------------------
; 2. denokv – waits for litestream to restore, then serves KV-Connect
; ---------------------------------------------------------------------------
[program:denokv]
command=/app/scripts/denokv_start.sh
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Give litestream time to restore before denokv opens the DB
startretries=5
startsecs=5
; denokv depends on litestream having done its initial restore
depends_on=litestream

; ---------------------------------------------------------------------------
; 3. Supercronic – runs the crontab for daily backup integrity checks
; ---------------------------------------------------------------------------
[program:supercronic]
command=/usr/local/bin/supercronic /app/crontab
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface
